# jboss_modified
#!/bin/sh
#
# /etc/init.d/jboss
# init script for Splunk.
# generated by 'splunk enable boot-start'.
#
# chkconfig: 2345 90 60
# description: Splunk indexer service
PIDFILE="/var/run/jboss/running.pid"
RETVAL=0
DAEMON_USER="jboss"
. /etc/init.d/functions
JBOSS_HOME=/home/jboss/jboss-5.1.0.GA
jboss_start() {
  pid="/var/run/jboss/running.pid"
  daemon --user $DAEMON_USER "nohup "/home/jboss/jboss-5.1.0.GA/bin/run.sh" -b 0.0.0.0 -c default > /dev/null 2> /dev/null &"
  sleep 6 
  name=$(ps -ef | grep -i java | grep jboss-5 | grep -v grep | awk -F " " '{print  $2}' > $pid)
  [ -x $name ] 
 RETVAL=$?
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/jboss
}
jboss_stop() {
  pid="/var/run/jboss/running.pid"
  name=$(cat $pid)
  var=$(kill -9 $name > /dev/null 2> /dev/null)
  [ -x $var > /dev/null ]
  RETVAL=$?
  [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/jboss
  
}
 
jboss_restart() {
jboss_stop
wait $$
jboss_start
}

jboss_status() {
  var=$(ps -ef | grep -i java | grep jboss-5 | awk -F " " '{print $2}')
  if [[ $var ]]
  then
  echo "Jboss is running:" $var
  else
  echo "Jboss is Not running"
  RETVAL=$?
  fi
}
 
jboss_startlog() {
var=$(ps -ef | grep -i java | grep jboss-5 | grep -v grep | awk -F " " '{print $2}')
 if [[ $var ]]
 then
  log=$(tail -100 $JBOSS_HOME/server/default/log/server.log | grep Started)
  echo $log
 else
  echo Currently Startlog not Available..
 fi

}
case "$1" in
  start)
    jboss_start
    ;;
  stop)
    jboss_stop
    ;;
  restart)
    jboss_restart
    ;;
  status)
   jboss_status
    ;;
  startlog)
   jboss_startlog
   ;;
  
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit $RETVAL
